# deploy/systemd/exam-study-automation.service
[Unit]
Description=Exam Study App Automation Loop
After=network.target

[Service]
Type=oneshot
User=jose
Group=jose
WorkingDirectory=/home/jose/exam-study-app
ExecStart=/bin/bash -lc '/home/jose/exam-study-app/scripts/run_automation.sh'
StandardOutput=append:/home/jose/exam-study-app/logs/automation.log
StandardError=append:/home/jose/exam-study-app/logs/automation.log
